# systemd/iterlife-app-deploy-webhook.service.d/10-log-perms.conf
[Service]
ExecStartPre=/usr/bin/install -d -o iterlife-reunion -g iterlife-reunion -m 2775 /apps/logs/webhook
